# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx_stdout.log
stderr_logfile=/var/log/supervisor/nginx_stderr.log
priority=10

[program:saleor]
command=/usr/local/bin/uvicorn saleor.asgi:application --host 127.0.0.1 --port 8080 --workers 2 --lifespan off --ws none --no-server-header --no-access-log --timeout-keep-alive 35 --timeout-graceful-shutdown 30 --limit-max-requests 10000
directory=/app
user=saleor
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/saleor_stdout.log
stderr_logfile=/var/log/supervisor/saleor_stderr.log
priority=20
environment=PYTHONUNBUFFERED="1"
